# etc/init.d/picoradiusd~
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

DEFAULT=/etc/default/radiusd
LOG_D=/var/log
RUN_D=/var/run
PID_F=$RUN_D/picorad.pid
RADACCT_D=/pico/db/acct
RADDB_D=/etc/picobill/raddb
IPADDR=$(ifconfig br-lan | sed -n 's/.*dr:\(.*\)Bc.*/\1/p')

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	mkdir -p $LOG_D
	mkdir -p $RUN_D
	mkdir -p $RADACCT_D
	mkdir -p $RADDB_D
	radiusd -i $IPADDR -p 1812,1813 -d $RADDB_D
}

stop() {
	[ -f $PID_F ] && kill $(cat $PID_F)
}

